# scripts/ansible/examples/proxy/etc/systemd/system/docker.service.d/http-proxy.conf
[Service]    
Environment="HTTP_PROXY=http://proxy.example.com:3128/" "HTTPS_PROXY=http://proxy.example.com:3129/" "NO_PROXY=localhost,127.0.0.1"

